# ping-after-wifi.service
[Unit]
Description=Ping target after Wifi is connected
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
Environment=INTERFACE=wlan0
Environment=PING_COUNT=4
ExecStart=/usr/local/bin/ping_target.sh

[Install]
WantedBy=multi-user.target
