# ping-after-wifi.service
[Unit]
Description=Ping target after Wi-Fi is connected
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/ping-after-wifi
ExecStart=/usr/local/bin/ping_target.sh

[Install]
WantedBy=multi-user.target
