# ping-after-wifi.service
[Unit]
Description=Ping target after WiFi is connected
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ping_target.sh
RemainAfterExit=yes
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
